# litellm-ollama-box.service
# litellm-ollama-box.service
# Host-side systemd USER service (not inside the container). Distrobox
# containers share your home directory, so hooking .bashrc for autostart
# would fire in every container you have, not just ollama-box. A systemd
# unit on the host that calls into the container avoids that.
#
# This starts ONLY the LiteLLM proxy inside ollama-box. It does NOT run
# llama-server (start-local-llama.sh) and does NOT load any model. The
# proxy config is local-only (no cloud entries, no ANTHROPIC_API_KEY needed
# anywhere), so leaving this running all the time costs nothing and bills
# nothing on its own. Whether Claude Code actually talks to it is
# controlled separately by claude-local-toggle.sh.
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp litellm-ollama-box.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now litellm-ollama-box.service
#
# To make this survive an actual machine boot (not just login), see
# README.md for the loginctl enable-linger step (install.sh also offers it
# as a prompt).
#
# Check status / logs:
#   systemctl --user status litellm-ollama-box.service
#   journalctl --user -u litellm-ollama-box.service -f
#
# When you actually want the local model available, start it manually
# (see README for what this launches with - KV cache quant, flash
# attention, speculative decoding flags), or use the "Start Local Model"
# desktop icon if you installed it:
#   ~/.local/bin/start-local-llama.sh

[Unit]
Description=LiteLLM proxy inside ollama-box distrobox (local-only, no cloud)
After=network.target

[Service]
Type=simple

# Edit the path to litellm_config.yaml to wherever you actually saved it.
# No ANTHROPIC_API_KEY needed, this config has no cloud model entries.
ExecStart=/usr/bin/distrobox enter ollama-box -- /bin/bash -lc "exec litellm --config /home/%u/litellm_config.yaml --port 4000"


# NOT "distrobox stop ollama-box" - that stops the whole container, which
# kills llama-server (started separately by start-local-llama.sh) as
# collateral damage every time this service restarts or stops. Only kill
# the litellm process itself; leave the container and anything else in it
# (the model server) running.
ExecStop=/usr/bin/distrobox enter ollama-box -- pkill -f "litellm --config"
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target


